# ansible/templates/nmc-update.timer.j2
[Unit]
Description=Daily NMC server update

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
