# ansible/templates/nmc-update.service.j2
[Unit]
Description=Update and rebuild NMC server

[Service]
Type=oneshot
ExecStart=/usr/local/bin/nmc-update
