# ansible/templates/nmc-server.service.j2
[Unit]
Description=NMC Server
{% if nmc_tracey_sidecar_enabled %}
After=network-online.target tracey-sidecar.service
Wants=network-online.target tracey-sidecar.service
{% else %}
After=network-online.target
Wants=network-online.target
{% endif %}

[Service]
User={{ nmc_run_user }}
Group={{ nmc_run_group }}
WorkingDirectory={{ nmc_build_dir }}
EnvironmentFile=-/etc/nmc/nmc.env
ExecStart={{ nmc_build_dir }}/nmc_server --port {{ nmc_port }} -k /etc/nmc/kubeconfig
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
